# buildroot/board/broadcom/iproc/skeleton/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -o remount,rw / # REMOUNT_ROOTFS_RW
::sysinit:/bin/mkdir -p /dev/pts
::sysinit:/bin/mkdir -p /dev/shm
::sysinit:/bin/mount -a
::sysinit:/bin/hostname -F /etc/hostname
# Try to mount custom file system on NAND
# ::sysinit:/usr/sbin/ubiattach -m 8 -d 1  # Try mtd8 first (spiflash used)
# ::sysinit:/usr/sbin/ubiattach -m 4 -d 1  # Then mtd4 if mtd8 not present
# ::sysinit:/bin/mkdir -p /mnt/custfs
# ::sysinit:/bin/mount -t ubifs /dev/ubi1_0 /mnt/custfs/
# ::sysinit:/bin/mount -t jffs2 /dev/mtdblock4 /mnt/custfs/ 
# now run any rc scripts
::sysinit:/etc/init.d/rcS

## initialize mdev
::sysinit:/sbin/mdev -s

## autorun the autorun script in /mnt/custfs
::sysinit:/mnt/custfs/autorun 

# Put a getty on the serial port
::respawn:/sbin/getty -L ttyS0 115200 vt100 # GENERIC_SERIAL

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/etc/init.d/rcK
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a

